# debian/postrm
#!/bin/sh -e

case "$1" in
    remove)
		systemctl disable sane-web-client
	;;

    purge)
        rm -rf /var/lib/sane-web-client
    ;;
esac

exit 0
