# deb/DEBIAN/prerm
#!/bin/sh -e

case "$1" in
    upgrade)
        systemctl stop scanpi
    ;;
    remove)
        #systemctl stop scanpi
        systemctl disable --now scanpi
    ;;
esac

exit 0
